# 3G_sakis3g_mod/connection.service
[Unit]
Description= 3G_connection

[Service]
User=root
ExecStart=/usr/bin/connection

[Install]
WantedBy=multi-user.target
